# linux/boot/rc.local
/bin/sh /home/pi/Scripts/boot/boot.sh
exit 0
